# tiger/.config/systemd/user/cros-garcon.service.d/override.conf
[Service]
Environment="QT_SCALE_FACTOR=0.8"
Environment="QT_QPA_PLATFORMTHEME=qt5ct"
